# root/etc/supervisor/supervisord.conf
[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; don't run supervisord as a daemon
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[program:crond]
user=root
command=/scripts/start/crond.sh -f
stderr_logfile = /var/log/supervisord/crond-stderr.log
stdout_logfile = /var/log/supervisord/crond-stdout.log

[program:apache]
user=root
command=/scripts/start/apache.sh -DFOREGROUND
stderr_logfile = /var/log/supervisord/apache-stderr.log
stdout_logfile = /var/log/supervisord/apache-stdout.log

[program:mysql]
user=mysql
command=/scripts/start/mysql.sh
stderr_logfile = /var/log/supervisord/mysql-stderr.log
stdout_logfile = /var/log/supervisord/mysql-stdout.log
